# nos-stock-daemon/etc/systemd/system/nos-stock-site.timer
[Unit]
Description=Update the P-Nos Point Site Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
